# packaging/systemd/shorewall.service
[Unit]
Description=Shorewall-nft firewall
Documentation=man:shorewall(8)
# Bring the firewall up before the network, so interfaces never come up
# unprotected. Do not run beside another firewall manager.
Wants=network-pre.target
Before=network-pre.target
Conflicts=firewalld.service ufw.service
After=

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/shorewall start
ExecStop=/usr/sbin/shorewall stop
ExecReload=/usr/sbin/shorewall reload
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
